# workspace/scripts/files/pioreactor.cron
0 0 */5 * * /usr/local/bin/pio run backup_database
0 0 */29 * * rm /home/pioreactorui/backend/build/static/exports/*.zip /home/pioreactorui/backend/build/static/exports/*.csv
