# workspace/scripts/files/system/systemd/conditionally_start_raspap.service
[Unit]
Description=Start local access point using RaspAP

[Service]
Type=simple
ExecStart=/boot/conditionally_start_raspap.sh

[Install]
WantedBy=multi-user.target
