# workspace/scripts/files/system/systemd/ngrok.service
[Unit]
Description=start ngrok on startup
After=network.target

[Service]
ExecStart=/opt/ngrok/ngrok start ui ws --config /home/pioreactor/.ngrok2/ngrok.yml
Type=simple

[Install]
WantedBy=multi-user.target
