# workspace/scripts/files/system/systemd/start_raspap.service
[Unit]
Description=Start local access point using RaspAP
Before=network-pre.target
Wants=network-pre.target

[Service]
ExecStart=/bin/bash /usr/local/bin/conditionally_start_raspap.sh
ExecStartPost=pio log -m "%n successful" -n systemd
Type=oneshot

[Install]
WantedBy=network.target
